# rootfs/etc/supervisor.d/mysql.ini

[program:mysql]
command        = /usr/bin/mysqld --user=mysql --console
stdout_logfile = /tmp/%(program_name)s.log
stderr_logfile = /tmp/%(program_name)s.log
autorestart    = true

